# engram-learn.timer
[Unit]
Description=The Crab Engram — Learn Timer (every 5 min)
Documentation=https://github.com/maisonnat/the-crab-engram

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
